# docker/samba/samba/etc/s6/nginx/run
#!/bin/sh

exec nginx -c '/var/lib/samba/nginx.conf' -g 'daemon off;'
