# ntp/input/rootfs/etc/init.d/data_prepare
#!/sbin/openrc-run
# shellcheck shell=ash

depend()
{
    need localmount
}

start()
{
    ebegin Preparing persistent data
    /sbin/ab_resizedata
    # make sure /data is mounted
    mount -a

    if [ ! -L /data/etc/resolv.conf ]
    then
        ln -sf /tmp/resolv.conf /data/etc/resolv.conf
    fi

    echo 'nameserver 10.1.0.53' > /tmp/resolv.conf

    mkdir -p /data/chrony && chown -R chrony:chrony /data/chrony

    eend 0
}
