# docker/nas/image/etc/s6/nmbd/run
#!/bin/sh

exec nmbd -F --no-process-group --log-stdout -s /etc/samba/smb.conf
