# docker/nas/samba/etc/s6/smbd/run
#!/bin/sh

printf 'guest\nguest\n' | smbpasswd -a guest

exec smbd -F --no-process-group --debug-stdout -s /etc/samba/smb.conf
