# docker/samba/image/etc/s6/winbindd/run
#!/bin/sh

exec winbindd -F --no-process-group -s /etc/samba/smb.conf
