# docker/samba/image/etc/s6/sshd/run
#!/bin/sh

mkdir -p /var/lib/samba/sshkeys
ssh-keygen -A

exec /usr/sbin/sshd -D
