# docker/samba/image/etc/s6/sshd/run
#!/bin/sh

mkdir -p /var/lib/samba/sshkeys
ssh-keygen -A
mv /etc/ssh/ssh_host_*_key* /var/lib/samba/sshkeys/

exec /usr/sbin/sshd -D
