# xmpp/prosody/s6/prosody/run
#!/bin/sh
exec s6-setuidgid prosody prosody -F
