# fadumper/backend/etc/s6/sshproxy/run
#!/bin/sh

set -e

exec s6-setuidgid '999:999' /usr/bin/ssh -i /config/icefox_ssh.key -D proxyuser@icefox.doridian.net
