# docker/samba/samba/etc/s6/smbd/run
#!/bin/sh

exec smbd -F --no-process-group -S -s /var/lib/samba/smb.conf
