# docker/nas/samba/etc/s6/nmbd/run
#!/bin/sh

exec nmbd -F --no-process-group --debug-stdout -s /etc/samba/smb.conf
