# docker/samba/image/etc/s6/smbd/run
#!/bin/sh

exec smbd -F --no-process-group --log-stdout -s /etc/samba/smb.conf
