# docker/samba/image/etc/s6/nmbd/run
#!/bin/sh

exec nmbd -F --no-process-group -s /etc/samba/smb.conf
